# add_bg.conf
start on startup
task
exec ~/DotFiles/add_bg.sh
